# deploy/opencodebot.service
[Unit]
Description=OpenCodez Telegram mirror bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bloob
Group=bloob
WorkingDirectory=/home/bloob/politia/projects/tg/opencodebot
Environment=OPENCODEBOT_CONFIG=/home/bloob/politia/state/projects/tg/opencodebot/config.json
EnvironmentFile=/home/bloob/politia/projects/tg/opencodebot/token.env
ExecStart=/usr/bin/node /home/bloob/politia/projects/tg/opencodebot/src/main.mjs
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
