# deploy/opencodebot.service
[Unit]
Description=OpenCodez Telegram mirror bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opencodebot
Group=opencodebot
WorkingDirectory=/opt/opencodebot
Environment=OPENCODEBOT_CONFIG=/etc/opencodebot/config.local.json
EnvironmentFile=/etc/opencodebot/token.env
ExecStart=/usr/bin/node /opt/opencodebot/src/main.mjs
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
